# back/stuff/supervisord.conf
 
[supervisord]

[program:websockets-serve]
command = python websocketserver.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 2
autorestart = true
